# config/launchd/io.github.app-store-monitor.daily.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.github.app-store-monitor.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/absolute/path/to/app-store-monitor/scripts/run_daily.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/absolute/path/to/app-store-monitor</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/absolute/path/to/app-store-monitor/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/absolute/path/to/app-store-monitor/logs/launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
